# deploy_electricity_costs/templates/electricity_costs.service.j2
[Unit]
Description="Service for {{ electricity_costs_user }}"

[Service]
ExecStart={{ electricity_costs_script }}
User={{ electricity_costs_user }}
